# lab3/1/etc/systemd/system/snap-snapd-18933.mount
[Unit]
Description=Mount unit for snapd, revision 18933
After=snapd.mounts-pre.target
Before=snapd.mounts.target
Before=local-fs.target

[Mount]
What=/var/lib/snapd/snaps/snapd_18933.snap
Where=/snap/snapd/18933
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
