# lab3/1/etc/udev/rules.d/70-snap.packetsender.rules
# This file is automatically generated.
# opengl
KERNEL=="dma_buf_te", TAG+="snap_packetsender_packetsender"
# opengl
KERNEL=="galcore", TAG+="snap_packetsender_packetsender"
# opengl
KERNEL=="mali[0-9]*", TAG+="snap_packetsender_packetsender"
# opengl
KERNEL=="nvhost-*", TAG+="snap_packetsender_packetsender"
# opengl
KERNEL=="nvmap", TAG+="snap_packetsender_packetsender"
# opengl
KERNEL=="pvr_sync", TAG+="snap_packetsender_packetsender"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_packetsender_packetsender"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_packetsender_packetsender"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_packetsender_packetsender"
# opengl
KERNEL=="vchiq", TAG+="snap_packetsender_packetsender"
# opengl
KERNEL=="vcsm-cma", TAG+="snap_packetsender_packetsender"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_packetsender_packetsender"
TAG=="snap_packetsender_packetsender", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_packetsender_packetsender $devpath $major:$minor"
